# docker/aio/cloudflared.conf
[program:cloudflared]
# Cloudflare Tunnel: an outbound-only connection to Cloudflare, which routes
# your public hostname to Caddy on loopback:80. No inbound ports, no
# port-forwarding, no static IP - the "run it from home" path. The tunnel is
# token-managed: create it in the Cloudflare Zero Trust dashboard and point the
# public hostname at http://localhost:80. The token is passed via TUNNEL_TOKEN
# (env, not argv, so it stays out of the process list).
command=/usr/local/bin/cloudflared tunnel --no-autoupdate run
user=sheaf
autorestart=true
startsecs=5
priority=25
environment=TUNNEL_TOKEN="%(ENV_CF_TUNNEL_TOKEN)s",TUNNEL_METRICS="127.0.0.1:0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
