# docker/aio/redis.conf
[program:redis]
# Transient state only (sessions, step-up challenges, rate-limit counters,
# presign cache). No persistence: a restart just re-authenticates users.
# Bound to loopback; only the backend in this container talks to it.
command=redis-server --bind 127.0.0.1 --port 6379 --save "" --appendonly no --maxmemory 256mb --maxmemory-policy allkeys-lru
user=sheaf
autorestart=true
startsecs=2
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
