# docker/aio/supervisord.conf
[supervisord]
nodaemon=true
user=root
pidfile=/run/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

# Bundled redis is dropped in here by the entrypoint only when REDIS_URL is
# not pointed at an external server.
[include]
files = /etc/supervisor/conf.d/*.conf

[program:app]
command=/usr/local/bin/aio-app-start.sh
user=sheaf
directory=/app
autorestart=true
startsecs=5
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caddy]
command=/usr/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
user=root
autorestart=true
startsecs=3
priority=30
# Persist ACME certs / state on the data volume so a container recreate does
# not re-fetch (and risk Let's Encrypt rate limits).
environment=XDG_DATA_HOME="/app/data",XDG_CONFIG_HOME="/app/data"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
